# server/services/worker/worker.service
[Unit]
Description=Have You Fed The Dog - worker service (push relay, IAP verify, cron)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/haveyoufedthedog/worker
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=PORT=3055
# PB superuser creds for the overdue cron (optional; see .env.example).
EnvironmentFile=-/opt/haveyoufedthedog/worker/.env

[Install]
WantedBy=multi-user.target
