# server/services/push-notifier/push-notifier.service
[Unit]
Description=Have You Fed The Dog — push notification relay
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/haveyoufedthedog/push-notifier
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=PORT=3055

[Install]
WantedBy=multi-user.target
